# systemd/audiocinema.service
[Unit]
Description=AudioCinema - Headless Analyzer Daemon
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%%PROJECT_DIR%%
ExecStart=%%PROJECT_DIR%%/.venv/bin/python -m src.headless_runner --daemon
Restart=always
RestartSec=5
# Hereda entorno de sistema; sin rutas absolutas de usuario

[Install]
WantedBy=multi-user.target
