# week4_v2_AudioCinemaApp/systemd/audiocinema.timer
[Unit]
Description=AudioCinema Analyzer (scheduled)

[Timer]
OnCalendar=__ONCALENDAR__
Persistent=true
Unit=audiocinema.service

[Install]
WantedBy=timers.target
